# images/php/supervisord/template/consumer.ini
[program:%event%]
directory=/var/www/html/apiapps
process_name=%event%
command=php artisan queue:work --queue=%event%
autostart=true
autorestart=true
startretries=999999
startsecs=5
redirect_stderr=true
