# images/php/supervisord/engine.ini
[program:jenkins]
directory=/home/search-engine/server
command=java -jar start.jar
process_name=solr
autostart=true
autorestart=true
startsecs=5
redirect_stderr=false
stdout_logfile=/home/salonese/salonese-backend/app/logs/search.out.log
stdout_capture_maxbytes=1MB
stderr_logfile=/home/salonese/salonese-backend/app/logs/search.error.log
stderr_capture_maxbytes=1MB